# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/www/storage/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:octane]
process_name=%(program_name)s
command=php /var/www/artisan octane:start --server=roadrunner --host=0.0.0.0 --port=8080 --watch
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/octane.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

[program:worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --sleep=3 --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

[program:scheduler]
process_name=%(program_name)s
command=php /var/www/artisan schedule:work
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=3600

[program:logs]
process_name=%(program_name)s
command=bash -c "tail -n 0 -F /var/www/storage/logs/octane.log | awk '{ print \"[OCTANE] \" $0; fflush() }' & tail -n 0 -F /var/www/storage/logs/worker.log | awk '{ print \"[WORKER] \" $0; fflush() }' & tail -n 0 -F /var/www/storage/logs/scheduler.log | awk '{ print \"[SCHEDULER] \" $0; fflush() }' & wait"
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stopwaitsecs=5
